# console/config/fotuka-imageprocessing-worker_dev.ini
[program:fotuka-imageprocessing-worker]
process_name=%(program_name)s_%(process_num)02d
command=/opt/local/bin/php yii image-processing
directory=/Users/rom/Sites/fotuka
autostart=true
autorestart=true
numprocs=4
stdout_logfile=/tmp/cron.log
redirect_stderr=true

